# deploy/agentd.service
[Unit]
Description=AgentOS Daemon
Documentation=https://github.com/agentos/agentos
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agentd
Group=agentd
WorkingDirectory=/var/lib/agentd
ExecStart=/usr/bin/python3 /usr/lib/agentd/agentd/daemon.py
Restart=on-failure
RestartSec=5s

# Environment
EnvironmentFile=-/etc/agentd/agentd.env
Environment=AGENTD_BASE=/var/lib/agentd
Environment=AGENTD_LOG_LEVEL=INFO

# Pi Zero memory safety
MemoryMax=256M
MemorySwapMax=0

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/agentd /var/log/agentd /tmp
ProtectHome=yes
CapabilityBoundingSet=
AmbientCapabilities=

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agentd

[Install]
WantedBy=multi-user.target
