# deploy/kriya.service
[Unit]
Description=Kriya Daemon
Documentation=https://github.com/agentos/agentos
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kriya
Group=kriya
WorkingDirectory=/var/lib/kriya
ExecStart=/usr/bin/python3 /usr/lib/kriya/kriya/daemon.py
Restart=on-failure
RestartSec=5s

# Environment
EnvironmentFile=-/etc/kriya/kriya.env
Environment=KRIYA_BASE=/var/lib/kriya
Environment=KRIYA_LOG_LEVEL=INFO

# Pi Zero memory safety
MemoryMax=256M
MemorySwapMax=0

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/kriya /var/log/kriya /tmp
ProtectHome=yes
CapabilityBoundingSet=
AmbientCapabilities=

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kriya

[Install]
WantedBy=multi-user.target
